# roles/dropbox/templates/dropboxd-supervisor.conf.j2
[program:dropboxd]
command={{ dropbox_install_dir }}/dropboxd
user={{ site_admin_user }}
startsecs=15